# systemd/hp-mute-led.service
[Unit]
Description=HP Mute LED Controller
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hp-mute-led.sh
Restart=always

[Install]
WantedBy=default.target
